# etc/xinetd.d/services
# entries for your /etc/services file.
qmserver	4242/tcp	# QMServer
qmclient	4243/tcp	# QMClient
