# usr/lib/systemd/system/scarletdmeserver@.service
[Unit]
Description=ScarletDME Per-Connection Telnet Server

[Service]
User=root
Group=qmusers
ExecStart=/usr/qmsys/bin/qm -n
StandardInput=socket

