# usr/lib/systemd/system/scarletdmeserver.socket
[Unit]
Description=ScarletDME Per-Connection Server

[Socket]
ListenStream=4242
Accept=true

[Install]
WantedBy=sockets.target

