# etc/systemd/system/scarletdmeclient@.service
[Unit]
Description=ScarletDME Per-Connection Telnet Server

[Service]
ExecStart=/usr/qmsys/bin/qm -n -q
StandardInput=socket

