# etc/systemd/system/scarletdmeclient.socket
[Unit]
Description=ScarletDME Per-Connection Telnet Server

[Socket]
ListenStream=4243
Accept=true

[Install]
WantedBy=sockets.target

