# etc/systemd/system/scarletdme@.service
[Unit]
Description=ScarletDME Per-Connection Server

[Service]
ExecStart=/usr/qmsys/bin/qm -n -q
StandardInput=socket

